# assets/ghost.conf
#/etc/init/ghost.conf
description "Ghost Blog"
author "Jeff Pihach <jpihach@gmail.com>"

setuid ubuntu

# Start the service after everything loaded
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

# Automatically restart service
respawn
respawn limit 99 5

chdir /var/www/ghost
env NODE_ENV=production

exec npm start 2>&1 >> /var/log/ghost.log
